# config/capsicum-relay.service
[Unit]
Description=capsicum push notification relay
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/repos/capsicum-relay
ExecStart=/home/deploy/.rbenv/shims/bundle exec puma -C config/puma.rb
Restart=on-failure
RestartSec=5
Environment=RACK_ENV=production
# Sentry (capsicum-relay#10). 実 DSN は commit せず、本番環境では以下のように設定:
# Environment=SENTRY_DSN=https://...@sentry.io/<project_id>
# Environment=SENTRY_RELEASE=<git SHA、deploy 時に git rev-parse HEAD で自動検出される場合は不要>

[Install]
WantedBy=multi-user.target
